# caddy/caddy-ipv6-prefixes.timer
[Unit]
Description=Update caddy IPv6 prefixes regularly

[Timer]
OnBootSec=10s
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
